# etc/cron.d/nextcloud
# min   hour    day     month   weekday   user      command
*/5     *       *       *       *         www-data  /usr/local/bin/nextcloud-cron
